# debian/postrm
#! /bin/sh -e

# summary of how this script can be called:
#        * <postrm> `remove'
#        * <postrm> `purge'
#        * <old-postrm> `upgrade' <new-version>
#        * <new-postrm> `failed-upgrade' <old-version>
#        * <new-postrm> `abort-install'
#        * <new-postrm> `abort-install' <old-version>
#        * <new-postrm> `abort-upgrade' <old-version>
#        * <disappearer's-postrm> `disappear' <r>overwrit>r> <new-version>
# for details, see /usr/share/doc/packaging-manual/

get_config() {
	db_get moodle/webserver
	webserver="$RET"
}

case "$1" in
    purge)
	. /usr/share/debconf/confmodule
	db_version 2.0

	get_config

	rm -f /etc/${webserver}/conf.d/moodle
	if [ -e "/etc/init.d/${webserver}" ]; then
	    invoke-rc.d ${webserver} reload || true
	fi

	rm -rf /usr/share/moodle
	ucf --purge /etc/moodle/config.php
	ucf --purge /etc/cron.d/moodle
	rm -rf /etc/moodle
	rm -rf /var/lib/moodle
	rm -f /etc/cron.d/moodle
	;;

	remove)
	. /usr/share/debconf/confmodule
	db_version 2.0

	get_config

	rm -f /etc/${webserver}/conf.d/moodle
	if [ -e "/etc/init.d/${webserver}" ]; then
	    invoke-rc.d ${webserver} reload || true
	fi
	;;

       upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)

        ;;

    *)
        echo "postrm called with unknown argument \`$1'" >&2
        exit 0

esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

